# supervisord.conf
[supervisord]
nodaemon=true

[program:custom1]
command=/start.sh

[program:Argo]
command=/usr/share/argo-tunnel.sh
