# supervisord.conf
[supervisord]
nodaemon=true

[program:custom1]
command=/start.sh
